# volto.service
[Unit]
Description=Volto LLM Command Server
After=network.target

[Service]
User=volto
WorkingDirectory=/home/volto/volto
ExecStart=/usr/bin/python3 /home/volto/volto/server.py
Restart=on-failure
RestartSec=5

[Install]
WantedBy=multi-user.target